# data/systemd/pamac-cleancache.timer
[Unit]
Description=Monthly clean packages cache

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
