# data/systemd/pamac-daemon.service
[Unit]
Description=Pamac Daemon

[Service]
Type=dbus
BusName=org.manjaro.pamac.daemon
ExecStart=/usr/bin/pamac-daemon
