# data/systemd/pamac-cleancache.service.in
[Unit]
Description=Clean packages cache

[Service]
Type=oneshot
ExecStart=@bindir@/pamac clean --no-confirm
