# data/systemd/pamac-system.service
[Unit]
Description=Pamac System Session

[Service]
Type=dbus
BusName=org.manjaro.pamac.system
ExecStart=/usr/bin/pamac-system-daemon
