# data/systemd/pamac-mirrorlist.service
[Unit]
Description=Generate mirrorlist

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-mirrors -f5
ExecStartPost=/usr/bin/pacman -Syy