# data/systemd/pamac.service
[Unit]
Description=Pamac

[Service]
Type=dbus
BusName=org.manjaro.pamac
ExecStart=/usr/lib/python3.3/site-packages/pamac/pamac-daemon.py
