# data/systemd/pamac-mirrorlist.timer
[Unit]
Description=Monthly generate mirrorlist

[Timer]
OnCalendar=monthly
RandomizedDelaySec=12h
Persistent=true

[Install]
WantedBy=timers.target
