# data/systemd/pamac-system.service.in
[Unit]
Description=Pamac System Session

[Service]
Type=dbus
BusName=org.manjaro.pamac.system
ExecStart=@bindir@/pamac-system-daemon
