# data/systemd/pamac-cleancache.service
[Unit]
Description=Clean packages cache

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
