# data/systemd/pamac-cleancache.timer
[Unit]
Description=Monthly clean packages cache

[Timer]
OnCalendar=Sat *-*-1..7 15:00:00
Persistent=true

[Install]
WantedBy=timers.target
